# music-atlas.service
[Unit]
Description=Music Atlas sync API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/music.193vbr.com
ExecStart=/usr/bin/python3 /srv/music.193vbr.com/server.py --db /var/lib/music-atlas/atlas.sqlite3 --host 127.0.0.1 --port 8090
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
