# music-atlas-backup.service
[Unit]
Description=Back up Music Atlas SQLite database

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/music-atlas-backup
