# music-atlas-backup.timer
[Unit]
Description=Daily Music Atlas database backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
