# deploy/electronic-genre-atlas-backup.timer
[Unit]
Description=Daily Electronic Genre Atlas database backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
